# crontab/crontab_remove_kq_data.sh
#!/bin/bash

/bin/rm -rf /netapp/show/pipe/works/kq_data/*